# rootfs/etc/services.d/vmware-task-telegram-bot/run
#!/usr/bin/with-contenv sh
vmware_task_bot -c /opt/vmware-task-telegram-bot/config.yml
